# supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:mynest-core]
command=/app/mynest
directory=/app
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/logs/mynest-core.log
stderr_logfile=/app/logs/mynest-core.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3

[program:telegram-bot]
command=/app/telegram-bot
directory=/app
environment=PLUGIN_MODE=grpc,PLUGIN_PORT=50051
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/logs/telegram-bot.log
stderr_logfile=/app/logs/telegram-bot.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/logs/nginx.log
stderr_logfile=/app/logs/nginx.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3