# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:mynest-core]
command=/app/mynest
directory=/app
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/mynest-core.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/nginx.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3